# dot_config/autostart/DankMaterialShell.desktop
[Desktop Entry]
Name=Dank Material Shell
Comment=DMS for Niri
Exec=dms run
Type=Application
Terminal=false
