# dot_config/autostart/niriusd.desktop
[Desktop Entry]
Name=Niri Utils Daemon
Comment=Helper daemon for Niri
Exec=niriusd
Type=Application
Terminal=false
