# dot_config/autostart/shikane.desktop
[Desktop Entry]
Name=Shikane
Comment=Background helper
Exec=shikane
Type=Application
Terminal=false
